# expiry-checker/self_hosted/cert-exp-check.cron
SHELL=/bin/BASH
MAILTO=certificates@mydomain.example
00 3 * * * /path/to/cert-exp-check.sh &> /path/to/log